# dbprod_slave/crontab/crontab
0 12 * * * root /scripts/backup_data.sh >> /var/log/cron.log 2>&1


